# packages/gemfire/profile.d/01-set-environ.sh
# Runs AFTER setup.d (Maintains Defintions)

PATH=$VMWARE_GEMFIRE_TARGET/bin:$PATH
